# deploy/systemd/pandora-central-monitor.service
[Unit]
Description=Pandora Central Monitoring Dashboard (Port 5000)
Documentation=https://github.com/yourusername/threat_project
After=network.target pandora-backend-admin.service

[Service]
Type=simple
User=pandora
Group=pandora
WorkingDirectory=/opt/pandora/central-monitor

# Environment
Environment="PYTHONPATH=/opt/pandora"
Environment="PYTHONUNBUFFERED=1"

# Start command (Flask)
ExecStart=/opt/pandora/central-monitor/venv/bin/python monitor_server.py

# Restart policy
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Resource limits
LimitNOFILE=4096
MemoryLimit=256M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pandora-monitor

[Install]
WantedBy=multi-user.target
