# deploy/systemd/pandora-backend-user.service
[Unit]
Description=Pandora Backend User API
After=network.target network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service redis.service elasticsearch.service

[Service]
Type=simple
User=pandora
Group=pandora
WorkingDirectory=/home/pandora/projects/pandora-threat-project/backend-user
ExecStart=/usr/bin/python3 -m uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONPATH=/home/pandora/projects/pandora-threat-project/backend-user"
Environment="PYTHONUNBUFFERED=1"

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

