# deploy/systemd/pandora-https-443.service
[Unit]
Description=Pandora HTTPS Server (Port 443) - Honeypot Listener
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/pandora/projects/pandora-threat-project/custom-webserver
ExecStart=/usr/bin/python3 /home/pandora/projects/pandora-threat-project/custom-webserver/port_443.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

