# deploy/systemd/pandora-backend-user.service
[Unit]
Description=Pandora User Backend API (Port 8000)
Documentation=https://github.com/yourusername/threat_project
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=pandora
Group=pandora
WorkingDirectory=/opt/pandora/backend-user

# Environment
Environment="PYTHONPATH=/opt/pandora"
Environment="PYTHONUNBUFFERED=1"

# Start command
ExecStart=/opt/pandora/backend-user/venv/bin/uvicorn api.main:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 2 \
    --log-level info

# Restart policy
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Resource limits
LimitNOFILE=65535
MemoryLimit=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pandora-backend-user

[Install]
WantedBy=multi-user.target
