# deploy/systemd/pandora-ids.service
[Unit]
Description=Pandora IDS Engine (Network Intrusion Detection)
Documentation=https://github.com/yourusername/threat_project
After=network.target postgresql.service

[Service]
Type=simple

# IDS cần chạy với quyền root để capture packets
User=root
Group=root

WorkingDirectory=/opt/pandora/ids

# Environment
Environment="PYTHONPATH=/opt/pandora"
Environment="PYTHONUNBUFFERED=1"

# Start command
ExecStart=/opt/pandora/ids/venv/bin/python ids_engine.py

# Restart policy
Restart=always
RestartSec=10

# Security (minimal vì cần root)
NoNewPrivileges=true
ProtectHome=true

# Resource limits (IDS có thể dùng nhiều CPU)
LimitNOFILE=65535
MemoryLimit=1G
CPUQuota=50%  # Giới hạn 50% CPU để không ảnh hưởng services khác

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pandora-ids

[Install]
WantedBy=multi-user.target
