# deploy/systemd/pandora-backend-admin.service
[Unit]
Description=Pandora Backend Admin API
After=network.target network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=pandora
Group=pandora
WorkingDirectory=/home/pandora/projects/pandora-threat-project/backend-admin
ExecStart=/usr/bin/python3 -m uvicorn api.main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONPATH=/home/pandora/projects/pandora-threat-project/backend-admin"
Environment="PYTHONUNBUFFERED=1"

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

